# mail-kanban-assistant/app/scheduler/launchd/com.local.mailassistant.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.local.mailassistant</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>mail-assistant</string>
      <string>run-daily</string>
      <string>--digest-out</string>
      <string>/Users/YOU/Library/Logs/mail-assistant/digest.md</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABS/PATH/TO/mail-kanban-assistant</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOU/Library/Logs/mail-assistant/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/Library/Logs/mail-assistant/stderr.log</string>
  </dict>
</plist>
